# configs/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock   ; (the path to the socket file)
chmod=0700                       ; sockef file mode (default 0700)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock ; use a unix:// URL  for a unix socket

[supervisord]
user            = root
nodaemon        = true
logfile         = /dev/null
logfile_maxbytes= 0

[program:rsyslog]
command         = rsyslogd -n
autostart       = true
autorestart     = true
startsecs       = 2
stopwaitsecs    = 2
stdout_logfile  = /dev/stdout
stderr_logfile  = /dev/stderr
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0

[program:postfix]
command         = /scripts/postfix.sh
autostart       = true
autorestart     = false
directory       = /etc/postfix
startsecs       = 0

[program:postsrsd]
command         = /scripts/postsrsd.sh
autostart       = true
autorestart     = true
startsecs       = 0

[program:dovecot]
command         = /scripts/dovecot.sh
autostart       = true
autorestart     = true
startsecs       = 0

[program:cron]
command         = /scripts/cron.sh
autostart       = true
autorestart     = true
startsecs       = 0

[program:opendkim]
command         = /scripts/opendkim.sh
user            = opendkim
autostart       = true
autorestart     = true
startsecs       = 5
stopwaitsecs    = 5
stdout_logfile  = /dev/stdout
stderr_logfile  = /dev/stderr
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0
